# macOS/crontab
00 6 * * * /Users/toor/ultragrid-nightly-mac.sh
30 6 * * * /Users/toor/ultragrid-nightly-mac-release.sh
0 16 * * * /Users/toor/ultragrid-nightly-mac.sh devel
