# macOS/crontab
00 7 * * * /Users/toor/ultragrid-nightly-mac.sh
30 7 * * * /Users/toor/ultragrid-nightly-mac-release.sh
0 16 * * * /Users/toor/ultragrid-nightly-mac.sh devel
